# airootfs/etc/systemd/serial-getty@ttyS0.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- \\u' --noclear --autologin suser --keep-baud 115200,57600,38400,9600 - $TERM
